# voider/call_caller
#!/sbin/openrc-run

start(){
/home/dollner/call_caller.sh
}
